# provisioning/virtual-machines/ansible/roles/gns3/files/gns3.service
[Unit]
Description=GNS3 Server

[Service]
ExecStart=/usr/bin/gns3server

[Install]
WantedBy=multi-user.target